# roles/debian/files/etc/init.d/sshd
#!/sbin/openrc-run

description="compiled OpenBSD Secure Shell server"
command="/opt/openssh/sbin/sshd"
command_background="YES"
pidfile="/run/sshd.pid"

log_it() {
    # if no supported loggers are found, fall back to a non-dependency solution
    f_logger=/usr/local/bin/logger
    if [ -x "$f_logger" ]; then
        "$f_logger" "$RC_SVCNAME" "$*"
    else
        echo "$(date +"%Y-%m-%d %H:%M:%S") ${RC_SVCNAME}: $*" >>/var/log/messages
    fi
}

check_privsep_dir() {
    # Create the PrivSep empty dir if necessary
    _d=/var/lib/sshd
    if [ ! -d "$_d" ]; then
        mkdir -p "$_d"
        chmod 0755 "$_d"
        log_it "Created PrivSep dir: $_d"
    fi
}

check_config() {
    /opt/openssh/sbin/sshd  -t || {
        msg="sshd config error"
        eerror "$msg"
    }
}

start_pre() {
    check_privsep_dir
    check_config
}

start_post() {
    log_it started
}

stop_post() {
    log_it stopped
}
